# templates/upstart.wikimetrics.conf.erb
start on <%= @service_start_on %>

setuid <%= scope.lookupvar('wikimetrics::user') %>
setgid <%= scope.lookupvar('wikimetrics::group') %>

script
    wikimetrics --mode <%= @mode %> --db-config <%= @config_directory %>/db_config.yaml --queue-config <%= @config_directory %>/queue_config.yaml --web-config <%= @config_directory %>/web_config.yaml
end script
